# nginx-hhvm/supervisord.conf
[program:hhvm]
user=www-data
command=hhvm --mode server --config /etc/hhvm/php.ini --config /etc/hhvm/server.ini

[program:nginx]
command=nginx -g "daemon off;"
