# nginx-symfony/supervisord.conf
[program:fpm]
command=php5-fpm -R -F

[program:nginx]
command=nginx -g "daemon off;"
